# systemd/mpd.service
[Service]
User=rodfer
ExecStart=/usr/bin/mpd
